# gnome-screensaver.service
[Unit]
Description=GNOME Screensaver

[Service]
BusName=org.gnome.ScreenSaver
ExecStart=@/usr/bin/gnome-screensaver gnome-screensaver --no-daemon
