# numlock.service
[Unit]
Description=numlockx

[Service]
ExecStart=/usr/bin/numlockx on

[Install]
WantedBy=desktop.target
