# set-background.service
[Unit]
Description=Set Desktop Background

[Service]
ExecStart=/home/matt/src/set_bg
