# xcompmgr.service
[Unit]
Description=X Compositing Manager

[Service]
ExecStart=@/usr/bin/xcompmgr xcompmgr
